# root/etc/supervisor/conf.d/mariadb.ini
[program:mariadb]
command=mysqld
user=mysql
redirect_stderr=true
stdout_logfile=/var/log/mysql.log
autorestart=true